# ansible/cockpit/templates/cockpit.socket.j2
[Unit]
Description=Cockpit Web Service Socket
Documentation=man:cockpit-ws(8)

[Socket]
ListenStream={{ cockpit_port }}

[Install]
WantedBy=sockets.target